# tools/curator/crontab
# Snapshot Every 1 hour
* * * * * curator /actions/snapshot.yml

